# templates/taskwarrior-sync.timer.j2
[Install]
WantedBy=default.target

[Unit]
Description=Runs taskwarrior synchronization periodically.

[Timer]
OnCalendar=*:0/1
Unit=taskwarrior-sync.service
